# download_from_moex_iss_intraday/crontab.cron
# m h  dom mon dow   user command
00 12 * * 1-5 date >> /market_data/log
00 12 * * 1-5 date >> /market_data/logerr

01 12 * * 1-5 python3 /load_data_from_moex_iss.py SBER >> /market_data/log 2>> /market_data/logerr
02 12 * * 1-5 python3 /load_data_from_moex_iss.py SBERP >> /market_data/log 2>> /market_data/logerr
